# generators/vulcanize/templates/sphinx/config/rubber/role/sphinx/crontab
<%
  @read_cmd = 'crontab -l'
  @write_cmd = 'crontab -'
  @additive = ["# rubber-sphinx-start", "# rubber-sphinx-end"]
%>

*/15 * * * * cd <%= RAILS_ROOT %> && rake ts:index
