# templates/graylog/config/rubber/role/graylog_server/graylog_server-upstart.conf
<%
  @path = "/etc/init/graylog-server.conf"
  @backup = false

  is_old_ubuntu = rubber_instance.os_version == '14.04'
  @skip = ! is_old_ubuntu
%>
description "graylog server"

start on (started mongodb and started elasticsearch)
stop on runlevel [016]

script
  cd <%= rubber_env.graylog_dir %>
  ulimit -n 65536
  exec java -jar graylog2-server.jar > graylog.log 2>&1
end script

post-start script
  status graylog-server | head -n1 | awk '{print $NF}' > <%= rubber_env.graylog_pid_file %>
end script

post-stop script
  rm -f <%= rubber_env.graylog_pid_file %>
end script
